# files/etc/uci-defaults/99-landscape
#!/bin/sh
/etc/init.d/landscape enable
exit 0
